# roles/calc-api-sum/templates/goapp.service.j2
[Unit]
Description=goapp

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/srv/calc-api-sum/goapp

[Install]
WantedBy=multi-user.target
